# config/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:script]
command=/bin/bash -c "/var/www/html/run.sh"
autorestart=true